# deploy/systemd/jarad-backend.service
[Unit]
Description=Jarad Backend
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=__BACKEND_USER__
Group=__BACKEND_USER__
WorkingDirectory=/opt/jarad/app
EnvironmentFile=-__REMOTE_ROOT__/backend/.env
Environment=JARAD_BACKEND_HOST=127.0.0.1
Environment=JARAD_BACKEND_PORT=8443
Environment=JARAD_DOCKER_HELPER=/usr/local/sbin/jarad-docker
Environment=JARAD_DB_PATH=/var/lib/jarad/jarad.sqlite3
ExecStart=/usr/bin/env JARAD_DB_PATH=/var/lib/jarad/jarad.sqlite3 JARAD_DOCKER_HELPER=/usr/local/sbin/jarad-docker /opt/jarad/app/.venv/bin/python -m uvicorn jarad_backend.main:app --host ${JARAD_BACKEND_HOST} --port ${JARAD_BACKEND_PORT}
Restart=on-failure
RestartSec=5
UMask=0077
# DNS access enforcement uses sudoers to allow only the narrow jarad-dns-access helper.
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
BindReadOnlyPaths=__REMOTE_ROOT__/backend:/opt/jarad/app
StateDirectory=jarad
StateDirectoryMode=0700
ReadWritePaths=/var/lib/jarad
PrivateDevices=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
